# templates/default/systemd.service.erb
# Generated by chef - Do not edit! <%# This .erb template is handmade. The 'do not edit' comment is intended for the resulting .conf :) %>

[Unit]
Description=RubDub XMPP PubSub Node

[Service]
User=<%= @service_user %>
WorkingDirectory=<%= @app_root %>
ExecStart=/usr/bin/nodejs <%= @app_root %>/src/index.js \
	-p <%= @port %> \
	-b <%= @bind_address %> \
	-d <%= @domain %> \
	-k <%= @tls_key_path %> \
	-c <%= @tls_cert_path %> \
	-r <%= @rollbar_key %> \
	>> <%= @log_path %> 2>&1
Restart=always

[Install]
WantedBy=multi-user.target

