# templates/default/upstart.conf.erb
#!upstart
# Generated by chef - Do not edit! <%# This .erb template is handmade. The 'do not edit' comment is intended for the resulting .conf :) %>

description "RubDub XMPP PubSub Node"

start on runlevel [2345]
stop on runlevel [06]

respawn
respawn limit 20 5
setuid <%= @service_user %>

script
	exec /usr/local/bin/node <%= @app_root %>/src/index.js \
	-p <%= @port %> \
	-b <%= @bind_address %> \
	-d <%= @domain %> \
	-k <%= @tls_key_path %> \
	-c <%= @tls_cert_path %> \
	>> <%= @log_path %> 2>&1
end script